# cron/crontab
* * * * * docker exec php /usr/local/bin/ips_task.sh >> /var/log/cron.log 2>&1
