# cron/crontab
* * * * * /usr/local/bin/run_with_logging.sh "ips task" "docker exec php /usr/local/bin/ips_task.sh"
0 3 * * * /usr/local/bin/run_with_logging.sh "logrotate" "docker container start -a logrotate"
